# deploy/host-caddy/systemd/caddy.service.d/override.conf
[Service]
ExecReload=
ExecReload=/usr/bin/caddy reload --address unix//var/lib/caddy/caddy-admin.sock --config /etc/caddy/Caddyfile --adapter caddyfile --force
